# computer/files/usr/lib/systemd/user/cli-update.timer
[Unit]
Description=Pull cli container image

[Timer]
OnBootSec=1h
OnUnitInactiveSec=1d

[Install]
WantedBy=timers.target
